# scripts/artifact.service
[Unit]
Description=ARTIFACT Arcade Machine
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kirniy/modular-arcade
Environment=ARTIFACT_ENV=hardware
Environment=PYTHONPATH=/home/kirniy/modular-arcade/src
# Let pygame auto-detect kmsdrm - do NOT set SDL_VIDEODRIVER explicitly!
ExecStart=/home/kirniy/modular-arcade/.venv/bin/python -m artifact.main
Restart=on-failure
RestartSec=5
# Give time for display to be available
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
