# scripts/vnvnc-video-wall.service
[Unit]
Description=VNVNC CRT Video Wall
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/kirniy/modular-arcade
Environment=PYTHONPATH=/home/kirniy/modular-arcade/src
Environment=ARCADE_DATA_DIR=/home/kirniy/modular-arcade/data
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=VNVNC_VIDEO_WALL_DEVICE=/dev/video0
Environment=VNVNC_VIDEO_WALL_CAPTURE_WIDTH=640
Environment=VNVNC_VIDEO_WALL_CAPTURE_HEIGHT=480
Environment=VNVNC_VIDEO_WALL_CAPTURE_FPS=25
Environment=VNVNC_VIDEO_WALL_OUTPUT_WIDTH=640
Environment=VNVNC_VIDEO_WALL_OUTPUT_HEIGHT=480
Environment=VNVNC_VIDEO_WALL_OUTPUT=pykms
Environment=VNVNC_VIDEO_WALL_DISPLAY_INDEX=1
Environment=VNVNC_VIDEO_WALL_DRM_CARD=/dev/dri/card1
Environment=VNVNC_VIDEO_WALL_DRM_CONNECTOR=HDMI-A-2
Environment=VNVNC_VIDEO_WALL_FBDEV=/dev/fb0
Environment=VNVNC_VIDEO_WALL_SWITCH_INTERVAL=0
Environment=VNVNC_VIDEO_WALL_PRIMARY_WINDOW=0
Environment=VNVNC_VIDEO_WALL_DISPLAY_YEAR=2017
ExecStart=/home/kirniy/modular-arcade/.venv/bin/python -m artifact.video_wall.renderer
Restart=no

[Install]
WantedBy=multi-user.target
